# configs/autostart/apollo-setup-io.service
[Unit]
Description=Apollo PipeWire I/O Device Setup
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/apollo-setup-io
Restart=on-failure
RestartSec=10
# Retry up to 3 times if Apollo isn't on PipeWire yet
StartLimitBurst=3
StartLimitIntervalSec=60

[Install]
WantedBy=default.target
